# contrib/systemd/ethllama.service
[Unit]
Description=ethicallama — ethical, local-only LLM inference server
Documentation=https://github.com/luluthehungrycat/ethicallama
After=network.target

[Service]
Type=simple
# Run as the same user that installed ethicallama. Edit this for your setup.
# Common values:
#   - User=root                (system-wide `pip install`)
#   - User=ethllama            (dedicated service user with system-wide install)
#   - User=moritz              (per-user `uv tool install` — adjust as needed)
User=ethllama
Group=ethllama
WorkingDirectory=/var/lib/ethllama
# Include both system-wide and per-user binary paths so the service works
# regardless of whether you used `pip install` or `uv tool install`.
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/home/ethllama/.local/bin:/home/moritz/.local/bin"
Environment="ETHLLAMA_CONFIG=/etc/ethllama/config.yaml"
EnvironmentFile=-/etc/ethllama/ethllama.env
ExecStart=ethllama serve --host 0.0.0.0 --port 10434
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536
# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ethllama
AmbientCapabilities=
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
